# api/container/etc/service/server/run
#!/bin/sh
echo /server/dist/index.js | entr -rn node /server/dist/index.js