# api/container/etc/service/server/run
#!/bin/sh
echo /server/bundled/index.js | entr -rnz node --enable-source-maps /_
