# api/container/etc/service/server/run
#!/bin/bash
cd /server/
if [[ ! -z "${DEVELOPMENT}" ]]; then
    npm run nodemon
else
    npm run serve
fi